# shopware/docker/build/supercronic-supervisord.ini
[program:supercronic]
command=/usr/local/bin/supercronic -quiet /home/www-data/crontab.txt
user=www-data
autostart=true
autorestart=true
stopwaitsecs=0
stopsignal=SIGTERM
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
